# files/system/etc/systemd/user/user-flatpak-setup.service.d/override.conf
[Service]
ExecStartPost=/usr/bin/systemctl disable --global --now user-flatpak-setup.timer