# files/systemd/system/skel-init.service
[Unit]
After=systemd-user-sessions.service
Before=greetd.service
ConditionPathExists=!/var/lib/dank-bluebuild/skel-init/run.completed

[Service]
Type=oneshot
ExecStart=/usr/libexec/dank-bluebuild/skel-init/run
ExecStartPost=/usr/bin/touch /var/lib/dank-bluebuild/skel-init/run.completed

[Install]
WantedBy=multi-user.target      