# files/systemd/user/dms.service.d/override.conf
[Unit]
PartOf=hyprland-session.target
After=hyprland-session.target
Requisite=hyprland-session.target

[Install]
WantedBy=hyprland-session.target