# files/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland Session Target
Requires=graphical-session.target
After=graphical-session.target
Wants=xdg-desktop-autostart.target