# files/system/etc/xdg/autostart/post-login-setup.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Type=Application
Name=Setup
Exec=/usr/libexec/dank-bluebuild/post-login-setup/auto
Icon=utilities-terminal
Terminal=false
Categories=Utility;